# deploy/systemd/trading-openclaw-backup.service
[Unit]
Description=Trading Agent - Back up OpenClaw state on config change

[Service]
Type=oneshot
ExecStart=/home/eswar/claude/trading-agent/tools/openclaw_backup.sh
